# unifi/root/etc/s6-overlay/s6-rc.d/20-config/up
#!/command/execlineb -P

with-contenv
fdmove -c 2 1
pipeline { /sbin/unifi_init/20-config.sh }
mawk -Winteractive "\{print \"[init-20-config] \" $0\}"