# armhf-mongodb/root/etc/s6-overlay/s6-rc.d/haveged/run
#!/command/execlineb -P

fdmove -c 2 1
pipeline { haveged -F -v1 }
mawk -Winteractive "\{print \"[haveged] \" $0\}"