# unifi/root/etc/s6-overlay/s6-rc.d/logging-setup/up
foreground {
    if { mkdir -p /var/log/unifi }
    if { chown nobody:nogroup /var/log/unifi }
        chmod 02755 /var/log/unifi
    }

