# armhf-mongodb/root/etc/s6-overlay/s6-rc.d/10-users/up
#!/command/execlineb -P

with-contenv
fdmove -c 2 1
pipeline { /sbin/mongo_init/10-users.sh }
mawk -Winteractive "\{print \"[init-10-users] \" $0\}"
