# unifi/root/etc/s6-overlay/s6-rc.d/30-keygen/up
#!/command/execlineb -P

with-contenv
fdmove -c 2 1
pipeline { /sbin/unifi_init/30-keygen.sh }
mawk -Winteractive "\{print \"[init-30-keygen] \" $0\}"